# docker/etc/symbol-consumer-template.ini
[program:%symbol%-stops]
command=/srv/app/bin/console messenger:consume cron --sleep=0.33
environment=RUNNING_WORKER="symbol-dedicated",PROCESSED_ORDERS="Stop",PROCESSED_SYMBOL="%symbol%"
autostart=true
autorestart=true
startretries=999
stdout_logfile=/srv/app/var/log/bot-supervizord-out.log
